# deploy/scriptforge.service
[Unit]
Description=CodeIncant Backend
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=scriptforge
Group=scriptforge
WorkingDirectory=/opt/scriptforge/backend

# Environment — ANTHROPIC_API_KEY is only needed if you add a server-side fallback.
# With BYOK it can be omitted entirely.
EnvironmentFile=-/opt/scriptforge/.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/opt/scriptforge/venv/bin/uvicorn app:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --log-level warning

Restart=always
RestartSec=5

# Sandbox containers need Docker socket access
SupplementaryGroups=docker

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
